# data/thermald.conf
# thermald - thermal daemon
# Upstart configuration file
# Manages platform thermals

description	"thermal daemon"

start on started dbus
stop on stopping dbus

respawn

script
echo "Waiting for thermal zone and coretemp sysfs to be ready "
sleep 5
exec thermald --no-daemon --dbus-enable
end script
