# data/thermald.conf
# thermald - thermal daemon
# Upstart configuration file
# Manages platform thermals

description	"thermal daemon"

start on (local-filesystems
	  and started dbus
	)
stop on stopping dbus

expect fork
respawn

exec thermald --no-daemon
