# data/thermald.service.in
[Unit]
Description=Thermal Daemon Service

[Service]
Type=dbus
BusName=org.freedesktop.thermald
ExecStart=@sbindir@/thermald --no-daemon --dbus-enable

[Install]
WantedBy=multi-user.target
Alias=dbus-org.freedesktop.thermald.service

