# data/thermald.conf
# thermald - thermal daemon
# Upstart configuration file
# Manages platform thermals

description	"thermal daemon"

start on started dbus
stop on stopping dbus

respawn

exec thermald --no-daemon --dbus-enable
