# data/thermald.service.in
[Unit]
Description=Thermal Daemon Service
ConditionVirtualization=no

# A udev rule may create this file to disable thermald on certain hardware
ConditionPathExists=!/run/thermald/firmware-managed

Wants=systemd-udev-trigger.service systemd-udev-settle.service
After=systemd-udev-trigger.service systemd-udev-settle.service

[Service]
Type=dbus
SuccessExitStatus=1
BusName=org.freedesktop.thermald
ExecStart=@sbindir@/thermald --systemd --dbus-enable

[Install]
WantedBy=multi-user.target
Alias=dbus-org.freedesktop.thermald.service

