# data/thermald.conf
# thermald - thermal daemon
# Upstart configuration file
# Manages platform thermals

description	"thermal daemon"

start on runlevel [2345] and started dbus
stop on stopping dbus

#
# don't respawn on error
#
normal exit 1

respawn

#
# consider something wrong if respawned 10 times in 1 minute
#
respawn limit 10 60

exec thermald --no-daemon --dbus-enable
