# data/thermald.service.in
[Unit]
Description=Thermal Daemon Service
Documentation=man:thermald(1)
Documentation=https://github.com/01org/thermal_daemon/blob/master/thermal_daemon_usage.txt
After=syslog.target

[Service]
Type=dbus
BusName=org.freedesktop.thermald
ExecStart=/usr/bin/thermald --no-daemon --dbus-enable
StandardError=null

[Install]
WantedBy=multi-user.target
Alias=dbus-org.freedesktop.thermald.service

